# deployment/roles/gateways/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter

[Service]
User={{ ansible_user }}
ExecStart=/home/{{ ansible_user }}/prometheus/node_exporter/node_exporter

[Install]
WantedBy=default.target
